# systemd/wp-force-monitor.service
[Unit]
Description=Force PipeWire defaults to analog sink and monitor source
After=pipewire.service pipewire-pulse.service wireplumber.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=%h/bin/wp-force-monitor.sh
